# systemd/vjt-claude-aup-watchdog.service
[Unit]
Description=vjt-claude AUP / idle / turns watchdog + post-clear memory scrubber
Documentation=https://github.com/vjt/claude-ircbot
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/vjt/code/IRC/vjt-claude
ExecStart=/usr/bin/python3 -u /home/vjt/code/IRC/vjt-claude/aup_watchdog.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/vjt/code/IRC/vjt-claude/aup_watchdog.log
StandardError=inherit

[Install]
WantedBy=default.target
